# crontab/cron
* * * * * node /usr/app/build/src/index.js
# remember to end this file with an empty new line
